# lede/target/linux/hisilicon/base-files/etc/uci-defaults/99-mkfs-f2fs
#!/bin/sh

if [ -b /dev/mmcblk0p13 ]; then
    mkfs.f2fs -f /dev/mmcblk0p13
fi

exit 0